# pwsh/profile.d/aliases.ps1
Set-Alias -Name fe -Value Invoke-FuzzyEdit
Set-Alias -Name fcd -Value Invoke-FuzzyCd
Set-Alias -Name fh -Value Invoke-FuzzyHistory
Set-Alias -Name z -Value __zoxide_z -Option AllScope -Scope Global -Force
Set-Alias -Name zi -Value __zoxide_zi -Option AllScope -Scope Global -Force
Set-Alias -Name help -Value Help-Command
Set-Alias -Name gsync -Value sync-git-branch
Set-Alias -Name g -Value git -Option AllScope -Scope Global -Force
Set-Alias -Name rebase -Value rebase-interactive
Set-Alias -Name fix-merge -Value resolve-git-conflict
Set-Alias -Name rbfast -Value rebase-f
Set-Alias -Name su -Value admin
Set-Alias -Name vim -Value $EDITOR
Set-Alias -Name ep -Value Edit-Profile
Set-Alias -Name reload -Value reload-profile
Set-Alias -Name help -Value Show-Help
Set-Alias -Name details -Value fastfetch

